# simple-ftp-server.service
[Unit]
Description=Simple FTP Server
After=network.target

[Service]
# 修改为实际的用户和组
User=root
Group=root

# 修改为实际的可执行文件路径和工作目录
WorkingDirectory=/opt/simple-ftp-server
ExecStart=/opt/simple-ftp-server/simple-ftp-server

# 自动重启配置
Restart=always
RestartSec=5s

# 资源限制 (可选)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
